# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:db-migrate]
command=alembic upgrade head
directory=/app
autostart=true
autorestart=false
startsecs=0
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1"

[program:bars-ingest]
command=bash -c "sleep 5 && exec python3 -m quant_daily_bars.cli bars ingest --schedule %(ENV_INGEST_INTERVAL)s"
directory=/app
autostart=true
autorestart=true
startretries=999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1"

[program:quant-daily-bars-api]
command=bash -c "sleep 5 && exec python3 -m quant_daily_bars.api.app"
directory=/app
autostart=true
autorestart=true
startretries=999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1"

[program:bars-backfill-gaps]
command=bash -c "sleep 5 && exec python3 -m quant_daily_bars.cli bars backfill-gaps --schedule %(ENV_INGEST_INTERVAL)s"
directory=/app
autostart=false
autorestart=false
startretries=999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1"
